# restartFrankBot.timer
# /etc/systemd/system/restartFrankBot.timer
[Timer]
OnActiveSec=1h
OnUnitActiveSec=1h

[Install]
WantedBy=timer.target
