# restartFrankBot.timer
# /etc/systemd/system/restartFrankBot.timer
[Timer]
OnActiveSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timer.target
